# packages/@letsdoc/server/backend.service
[Unit]  
Description=Letsdoc backend service 
Wants = network-online.target
After = network.target network-online.target

[Service]      
ExecStart=/usr/local/go/bin/go run main.go   
WorkingDirectory=/home/ubuntu/actions-runner/_work/letsdoc/letsdoc/packages/@letsdoc/server

User=ubuntu
Group=ubuntu

Restart=always    
RestartSec=15
KillSignal=SIGINT

SyslogIdentifier=backend    
PrivateTmp=true  

[Install]      
WantedBy=multi-user.target
